# conf/files-backfire_ar71xx/files/etc/rc.local
#!/bin/sh

# /etc/init.d/pure-ftpd disable
# /etc/init.d/samba disable

exit 0
